# config/supervisord.conf
[supervisord]
#nodaemon=true
nocleanup = true

[program:fluentd]
command=/usr/local/rvm/bin/rvm-shell -l -c "fluentd -c ./fluent/fluent.conf"
user=fluentd
autostart=true
autorestart=true
